# Desktop/CoolRune-Files/etc/s6/current/run-image/service/s6-linux-init-early-getty/run
#!/usr/bin/execlineb -P

/usr/bin/agetty -L -8 tty7 115200
